# root/etc/services.d/YACReaderLibraryServer/run
#!/bin/execlineb -P

# Redirect stderr to stdout.
fdmove -c 2 1

# Drop privileges and set env
with-contenv s6-setuidgid abc

# Start YACReaderLibraryServer
YACReaderLibraryServer --loglevel debug start
